# deploy/claude-ide-bridge.service
[Unit]
Description=Claude IDE Bridge MCP Server
Documentation=https://github.com/Oolab-labs/claude-ide-bridge
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/claude-ide-bridge

# Load personal config — FIXED_TOKEN, PORT, WORKSPACE must be set here
EnvironmentFile=/root/claude-ide-bridge/.env.vps

# Bridge process — reads PORT, WORKSPACE, FIXED_TOKEN from EnvironmentFile
ExecStart=/usr/bin/node /root/claude-ide-bridge/dist/index.js \
    --port ${PORT} \
    --workspace ${WORKSPACE} \
    --fixed-token ${FIXED_TOKEN} \
    --grace-period 120000 \
    --vps

# Graceful shutdown: SIGTERM, wait 10s, then SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

# Restart policy: always restart unless manually stopped
Restart=always
RestartSec=5

# Logging — view with: journalctl -u claude-ide-bridge -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-ide-bridge

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
